# gateman/files/gateman.service
[Unit]
Description=Gateman Daemon

[Service]
ExecStart=/usr/local/bin/gateman

[Install]
WantedBy=multi-user.target
