# gateman/files/gateman.service
[Unit]
Description=Gateman

[Service]
ExecStart=/usr/local/bin/gateman --address 0.0.0.0 --port 9000

[Install]
WantedBy=multi-user.target
